# backend/supervisord.conf
[supervisord]
nodaemon=true
logfile=/dev/stdout
logfile_maxbytes=0
pidfile=/tmp/supervisord.pid

[program:api]
command=uvicorn autoflipr.main:app --host 0.0.0.0 --port 8000
directory=/app
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=10

[program:worker-scrape]
command=celery -A autoflipr.celery_app worker -Q scrape -c 2 -n worker-scrape@localhost --loglevel=info
directory=/app
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=20

[program:worker-scan]
command=celery -A autoflipr.celery_app worker -Q scan -c 2 -n worker-scan@localhost --loglevel=info
directory=/app
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=20

[program:worker-llm-0]
command=celery -A autoflipr.celery_app worker -Q llm -c 1 -n worker-llm-0@localhost --loglevel=info
directory=/app
environment=GEMINI_KEY_INDEX="0"
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=20

[program:worker-llm-1]
command=celery -A autoflipr.celery_app worker -Q llm -c 1 -n worker-llm-1@localhost --loglevel=info
directory=/app
environment=GEMINI_KEY_INDEX="1"
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=20

[program:worker-llm-2]
command=celery -A autoflipr.celery_app worker -Q llm -c 1 -n worker-llm-2@localhost --loglevel=info
directory=/app
environment=GEMINI_KEY_INDEX="2"
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=20

[program:worker-score]
command=celery -A autoflipr.celery_app worker -Q score,mot,default -c 4 --loglevel=info
directory=/app
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=20

[program:beat]
command=celery -A autoflipr.celery_app beat --loglevel=info --schedule=/tmp/celerybeat-schedule
directory=/app
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=20
